# multi-node-example/nfs-docker/etcd.conf
start on runlavel [2345]
stop on runlavel [!2345]

exec /bin/etcd --listen-client-urls 'http://127.0.0.1:4001' --data-dir '/data'
